# debian/pt-firmware-updater.service
[Unit]
Description=pi-top device firmware updater daemon

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/python3 -B /usr/bin/pt-firmware-checker --loop-time 3 --wait-time 300 --max-wait-timeout 3600

[Install]
WantedBy=graphical.target
